# data/com.unsplash.wallpaper.service
[Unit]
Description=Unsplash Wallpaper Manager - Automatic wallpaper changer
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/unsplash-wallpaper --run-job
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
